# initial/etc/local.d/00apk-update-upgrade.start
#!/bin/sh
echo "apk update & upgrade"

/sbin/apk update

OUTPUT=$(/sbin/apk upgrade)
echo "$OUTPUT"
echo "$OUTPUT" | grep -E "Installing|Upgrading" > /dev/null
if [ $? -eq 0 ]
then
  echo "Something was installed or upgraded - rebooting"
  /sbin/reboot
fi

echo "apk update & upgrade - done"
