# initial/etc/local.d/10dehydrated.start
#!/bin/sh

if [ -z "$(ls /etc/dehydrated/accounts)" ]; then
  echo "Registering letsencrypt account"
  /usr/local/bin/dehydrated --register --accept-terms
  /sbin/lbu commit
fi

echo "Requesting certificates"
/usr/local/bin/dehydrated -c
